# v2ex_spider/v2ex_crontab.cron
0 11 * * * /usr/src/v2ex/v2ex_spider/start_spider.sh
